# Wordpress Backup.cron
0 1 * * * /usr/bin/mysqldump --defaults-extra-file=/home/ubuntu/.quicklogin.cnf -u root --single-transaction --quick --lock-tables=false --all-databases > /home/ubuntu/full-backup-$(date +\%F).sql
0 1 * * * /bin/tar -czvf /home/ubuntu/html-backup-$(date +\%F).tar.gz /var/www/html
0 1 * * * su -c "/usr/local/bin/wp plugin update --path='/var/www/html/' --all" -s /bin/bash www-data
0 1 * * * su -c "/usr/local/bin/wp core update --path='/var/www/html/' " -s /bin/bash www-data
